# files/etc/uci-defaults/image-selector
autocore luci luci-i18n-base-fr luci-app-upnp luci-theme-argon luci-app-argon-config
